# files/zypper.cron
#! /usr/bin/env bash

/usr/bin/zypper refresh

